# polling-agent.service
[Unit]
Description=App Backend Polling Agent
After=network-online.target capture-daemon.service
Wants=network-online.target
Requires=capture-daemon.service

[Service]
Type=simple
User=autonexego30
WorkingDirectory=/home/autonexego30/Egocentric_Data_capture_Hardware_mono16
EnvironmentFile=/home/autonexego30/Egocentric_Data_capture_Hardware_mono16/.env.polling
ExecStart=/usr/bin/python3 /home/autonexego30/Egocentric_Data_capture_Hardware_mono16/polling_agent.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
